# trunk/system/config-templates/default-templates/mainserver/all/rootdir/etc/udev/rules.d/85-skolesys-usb-backup.rules
ACTION=="add", BUS=="usb", SUBSYSTEM=="block", RUN+="/usr/sbin/ss_udevparser"
ACTION=="remove", SUBSYSTEM=="block", RUN+="/usr/sbin/ss_udevparser"
